# debian/std-fpm.postinst
#!/bin/bash
systemctl enable std-fpm.service
systemctl start std-fpm
