# debian/std-fpm.postinst
#!/bin/bash
chown www-data:www-data /run/std-fpm/
systemctl enable std-fpm.service
systemctl start std-fpm
