# debian/std-fpm.postrm
#!/bin/bash
systemctl stop std-fpm
systemctl disable std-fpm.service
